# wss_server/supervisord.conf
[supervisord]
nodaemon=true

[program:Heart_Beat]
autostart=true
autorestart=false
startsecs=0
command=/bin/sh /usr/bin/start-heartbeat.sh


